# deploy/bettracker-worker.service
[Unit]
Description=BetTracker Scan Worker
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/bettracker
EnvironmentFile=/opt/bettracker/.env
ExecStart=/home/ubuntu/.local/bin/uv run python -m src.workers.scan_worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bettracker-worker
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/bettracker
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
UMask=0077

[Install]
WantedBy=multi-user.target
